# systemd/postmaster.socket
[Unit]
Description=PostMaster Socket

[Socket]
ListenStream=/run/postoffice/postmaster.socket

[Install]
WantedBy=sockets.target